# configs/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=0
loglevel=info

# ========================================
# PostgreSQL 数据库
# ========================================
[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres -D %(ENV_PGDATA)s -c max_connections=300 -c unix_socket_directories=/run/postgresql
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopwaitsecs=30
stopasgroup=true
killasgroup=true

# ========================================
# LiteLLM 网关
# ========================================
[program:litellm]
command=/app/scripts/entrypoint_allinone.sh litellm
user=root
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopwaitsecs=10
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED="1",DATABASE_URL="postgresql://llmproxy:dbpassword9090@127.0.0.1:5432/litellm"

# ========================================
# TrajProxy 主应用
# ========================================
[program:traj_proxy]
command=/app/scripts/entrypoint_allinone.sh traj_proxy
user=root
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopwaitsecs=15
stopasgroup=true
killasgroup=true
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/app",DATABASE_URL="postgresql://llmproxy:dbpassword9090@127.0.0.1:5432/traj_proxy"

# ========================================
# Nginx 反向代理
# ========================================
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
user=root
autostart=true
autorestart=true
priority=40
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
stopwaitsecs=5
